# root/etc/s6-overlay/s6-rc.d/init-config-obico/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# make folders
mkdir -p \
	/config/log/redis \
	/config/media \
	/run/obico

# permissions
chown -R abc:abc \
	/config \
	/app \
	/run/obico
