# root/etc/s6-overlay/s6-rc.d/init-config-obico/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# make folders
mkdir -p \
	/config/log/{redis,obico} \
	/config/media \
	/run/obico

if [ ! -f "/.setup" ]; then
	echo "Initialising obico-server (First time setup)"
	cd /app/obico/backend

	s6-setuidgid abc \
		python3.7 manage.py migrate >/dev/null &&
		python3.7 manage.py collectstatic -v 2 --noinput >/dev/null

	touch "/.setup"
fi

# permissions
chown -R abc:abc \
	/config \
	/app \
	/run/obico
