# root/etc/s6-overlay/s6-rc.d/svc-ml_api/run
#!/usr/bin/with-contenv bash

export DEBUG="True"
export FLASK_APP="server.py"

exec \
    cd /app/obico/ml_api s6-setuidgid abc \
        gunicorn --bind 0.0.0.0:3333 --workers 1 wsgi
