# root/etc/s6-overlay/s6-rc.d/svc-ml_api/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

cd /app/obico/ml_api

exec s6-setuidgid abc \
	gunicorn --bind 0.0.0.0:3333 --workers 1 wsgi >/config/log/obico/ml_api.log
