# root/etc/services.d/celery/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

cd /app/obico/backend

exec s6-setuidgid abc \
	celery -A config worker --beat -l info -c 2 -Q realtime,celery



