# root/etc/s6-overlay/s6-rc.d/svc-backend/run
#!/usr/bin/with-contenv bash

exec \
    cd /app/obico/backend s6-setuidgid abc \
        daphne -b 0.0.0.0 -p 3334 config.routing:application
