# root/etc/s6-overlay/s6-rc.d/svc-celery/run
#!/usr/bin/with-contenv bash

exec \
    cd /app/obico/backend s6-setuidgid abc \
        celery -A config worker --beat -l info -c 2 -Q realtime,celery
